# deploy/server-agent.service
[Unit]
Description=server-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/server_agent
EnvironmentFile=/etc/server-agent/agent.env
ExecStart=/usr/bin/npm start --workspace=@server-agent/server
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/server-agent/db
MemoryMax=512M

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
